# dssd.service
[Unit]
Description=Dead Simple Secret Daemon

[Service]
Type=dbus
BusName=org.freedesktop.secrets
ExecStart=/usr/bin/dssd
